# units/etc/sv/hepcli/run
#!/bin/bash

set -m
exec /usr/local/bin/heplify -hs ${HEP_SERVER}:${HEP_PORT}